# weather.timer
[Unit]
Description=Run weather hourly.

[Timer]
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
